# bluealsa-monitor.service.in
[Unit]
Description=Bluealsa PCM manager

[Service]
Type=simple
ExecStart=@prefix@/bin/bluealsa-monitor

[Install]
WantedBy=default.target
