# bluealsa-monitor.service.in
[Unit]
Description=Bluealsa PCM manager
Documentation=https://github.com/borine/bluealsa-monitor

[Service]
Type=simple
ExecStart=@bindir@/@monitor_script@

[Install]
WantedBy=default.target
